# provision/ansible/roles/nas/templates/lib/systemd/system/zfs-scrub@.service.j2
# {{ ansible_managed }}
[Unit]
Description=zpool scrub on %i

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=/usr/sbin/zpool scrub %i